# anbox-sensors.conf
# anbox-sensors

description "Anbox sensors"

start on started anbox-container

exec anbox-sensord
