# waydroid-sensors.conf
# waydroid-sensors

description "Waydroid sensors"

start on started waydroid-container

exec waydroid-sensord
